# dist/ban-failed.service
[Unit]
Description=Ban Failed

[Service]
ExecStart=/usr/sbin/ban-failed -f

[Install]
WantedBy=multi-user.target
